# systemd/pllm-rdma-state.service
[Unit]
Description=PLLM remote transactional live-state RDMA store
After=network-online.target

[Service]
Type=simple
WorkingDirectory=@ROOT_DIR@
ExecStart=@ROOT_DIR@/rdma_bridge/build/pllm-rdma-store --server /mnt/ssd-storage/pllm-cache/live-state --port 17901 --allocator cuda-host --token-file @HOME@/.config/pllm/rdma-token
Restart=on-failure
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=default.target
